# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:crucible]
command=sh /opt/fecru-3.5.4/bin/start.sh