# supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:crucible]
command=sh /opt/fecru/bin/start.sh